# systemd/wifi-dashboard.service
[Unit]
Description=WiFi Monitor — Dashboard Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pedrogariglio
WorkingDirectory=/home/pedrogariglio/wifi-stability-monitor
ExecStart=/home/pedrogariglio/wifi-stability-monitor/.venv/bin/uvicorn dashboard.main:app --host 0.0.0.0 --port 8088
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
