# systemd/wifi-metrics.timer
[Unit]
Description=WiFi Metrics Collector — cada 30 segundos

[Timer]
OnBootSec=60
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
